# deb-updater.service
[Unit]
Description=Debian Package Auto-Updater
After=network.target

[Service]
ExecStart=/usr/bin/deb-updater
Restart=on-failure

[Install]
WantedBy=multi-user.target
