# debian/deb-updater/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reexec || true
systemctl enable deb-updater.service || true
systemctl start deb-updater.service || true
